# tomat.service
[Unit]
Description=Tomat Pomodoro Timer Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=tomat daemon run
Restart=always
RestartSec=5

[Install]
WantedBy=default.target