# examples/systemd.service
# Systemd user service for tomat (manual installation)
#
# RECOMMENDED: Use `tomat daemon install` for automatic installation
#
# This file is for manual installation only. It assumes tomat is installed
# via `cargo install tomat` and available in ~/.cargo/bin/
# 
# Manual installation steps:
# 1. Copy this file to ~/.config/systemd/user/tomat.service
# 2. Adjust the ExecStart path if tomat is installed elsewhere
# 3. Run: systemctl --user daemon-reload
# 4. Run: systemctl --user enable tomat.service
# 5. Run: systemctl --user start tomat.service

[Unit]
Description=Tomat Pomodoro Timer Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/tomat daemon run
Restart=on-failure
RestartSec=5
# Inherit user's PATH for hooks to find system commands (e.g., notify-send)
Environment="PATH=/run/current-system/sw/bin:/etc/profiles/per-user/%u/bin:%h/.nix-profile/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target